# deploy/transcription-bot.service
[Unit]
Description=Telegram Transcription Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/transcription-bot
Environment="PATH=/home/ubuntu/transcription-bot/venv/bin"
ExecStart=/home/ubuntu/transcription-bot/venv/bin/python bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
